# systemd/user/gen-agenda.service
[Unit]
Description=Generate agenda

[Service]
Type=oneshot
ExecStart=/bin/bash /home/elkorn/Documents/Bash/genAgenda.sh