# systemd/user/update-dotfiles.timer
[Unit]
Description=Schedule dotfile updates

[Timer]
OnBootSec=15s
OnCalendar=*:*:0,15,30,45

[Install]
WantedBy=timers.target
