# systemd/user/update-dotfiles.timer
[Unit]
Description=Schedule dotfile updates

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
