# systemd/user/update-dotfiles.service
[Unit]
Description=Update dotfiles

[Service]
Type=oneshot
WorkingDirectory=/home/elkorn/.config/dotfiles/
ExecStart=/usr/bin/make all